# systemd/eos-btrfs-homesnapshot.timer
[Unit]
Description=Create a btrfs snapshot of the /home subvolume daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target