# systemd/eos-btrfs-prune-homesnapshot.timer
[Unit]
Description=Prune the btrfs snapshots of the /home subvolume daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target