# systemd/eos-btrfs-prune-rootsnapshot.service
[Unit]
Description=Prune snapshots of the / subvolume

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/prune-snapshot.sh /.snapshots/ 7