# systemd/eos-btrfs-prune-rootsnapshot.timer
[Unit]
Description=Prune the btrfs snapshots of the / subvolume daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target